# docker-mods/pm-base-tools/root/etc/services.d/log-viewer/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc /usr/bin/frontail -U ${FRONTAIL_USERNAME} -P ${FRONTAIL_PASSWORD}